# crontab
# Memory sync cron job - runs every 5 minutes to ensure perfect sync
# WORKDIR in Dockerfile is /app, Python packages are in /root/.local
*/5 * * * * cd /app && PYTHONPATH=/root/.local/lib/python3.11/site-packages python memory_sync_cron.py >> /var/log/memory_sync.log 2>&1
